# aa_wireless_dongle/package/aawg/src/aawgd.service
[Unit]
Description=Android Auto Wireless Gateway Daemon
After=bluetooth.service network.target
Wants=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m aawgd
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
